# systemd/gpu-guard.service.in
[Unit]
Description=GPUGuard - kill processes on GPUs flagged as stalled by GPUMonitor
After=network-online.target gpu-monitor.service
Wants=network-online.target
Wants=gpu-monitor.service

[Service]
Type=simple
# Runs as root so it can terminate any user's compute processes.
User=root
EnvironmentFile=@BASE@/config.env
Environment=TZ=UTC
ExecStart=/usr/bin/python3 @BASE@/bin/gpu_guard.py
Restart=always
RestartSec=5
# The killer must never be OOM-reaped before the jobs it polices.
OOMScoreAdjust=-500
# Hardening: the guard needs root + kill(2) on other users' procs, but nothing
# else. Shrink the blast radius of any compromise. (It must NOT set NoNewPrivileges
# or restrict SIGKILL to other users, so leave PID capabilities intact.)
ProtectSystem=strict
ReadWritePaths=@BASE@/data
ProtectHome=read-only
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
RestrictRealtime=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
